# gcp/supervisord.conf
[supervisord]
nodaemon=true
user=appuser
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/tmp
logfile_maxbytes=50MB
logfile_backups=3
loglevel=info

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ============================================================================
; Express Server (Node.js) - Main API Gateway
; ============================================================================
[program:express]
command=node --import tsx server/index.ts
directory=/app
user=appuser
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",PORT="5000",PYTHON_BACKEND_URL="http://127.0.0.1:8000"

; ============================================================================
; FastAPI Server (Python) - AI/ML Services
; ============================================================================
[program:fastapi]
command=python3.12 -m uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers %(ENV_WORKERS)s --loop uvloop --http httptools
directory=/app
user=appuser
autostart=true
autorestart=true
startsecs=15
startretries=3
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1"

; ============================================================================
; Process Group
; ============================================================================
[group:backend]
programs=express,fastapi
priority=999
